# babymeal.service
[Unit]
Description=치밀한 이유식 - Baby Meal Cube Manager
After=network.target mysqld.service

[Service]
Type=simple
User=daero52
WorkingDirectory=/volume1/DR_DATA1/babyMeal
ExecStart=/bin/python3 /volume1/DR_DATA1/babyMeal/web/app.py --port 8990
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
